# maintenance/systemd/backup-docker-volumes.timer
[Unit]
Description=Rsync docker volumes to a remote called "gdrive" that should be preconfigured

[Timer]
Unit=backup-docker-volumes.service
OnCalendar=daily
Peristent=True

[Install]
WantedBy=timers.target
